# slammer/user-sync/src/install/sync.cron
#
# Regular cron jobs for the SlammerSync
#

30 * * * * root [ -x #SYNC_BIN# ] && #SYNC_BIN# -s all -c all > /dev/null 2>&1
